# packages/xmedia/torrent/acephproxy/system.d/acephproxy.service
[Unit]
Description=AcePHProxy server
Requires=network-online.target
After=network-online.target acestream.service
Before=kodi.service
ConditionPathExists=/storage/.cache/services/acephproxy.conf

[Service]
ExecStartPre=-/usr/bin/acephp-config
ExecStart=/bin/sh -c '/usr/bin/php /storage/.config/acephproxy/acePHP.php 2>/var/log/acephproxy.log'
ExecStartPost=-/usr/bin/acephp-list.start
TimeoutStopSec=5
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
