# projects/S905/packages/32bit-libs/system.d/retroarch.service
[Unit]
Description=RetroArch emulator frontend
ConditionPathExists=/var/lock/start.retro

[Service]
Environment=DISPLAY=:0.0
Environment=SDL_MOUSE_RELATIVE=0
Environment=FB_MULTI_BUFFER=2
Environment=HOME=/storage
Environment=LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/lib32:/lib32:/usr/lib/pulseaudio:/usr/lib32/pulseaudio
ExecStart=/usr/bin/retroarch
KillMode=process
TimeoutStopSec=3
Restart=on-failure
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
